# themes/luci-theme-avalon/root/etc/uci-defaults/theme-avalon
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Avalon=/luci-static/avalon
	set luci.main.mediaurlbase=/luci-static/avalon
	commit luci
EOF
