# themes/sbadmin/root/etc/uci-defaults/luci-theme-sbadmin
#!/bin/sh
uci batch <<-EOF
	set luci.themes.SbAdmin=/luci-static/sbadmin
	set luci.main.mediaurlbase=/luci-static/sbadmin
        commit luci
EOF
	
